# etc/config/includes.chroot/etc/init/pine-hacks.conf
description "Pine64 hacks"

task

start on startup

script
    # unmute sound
    amixer -c 1 set 'AIF1 Slot 0 Digital DAC' unmute
end script
