# etc/config/includes.chroot/etc/init/set-info.conf
description "Update machine info"

task

start on started dbus

script
    hostnamectl --pretty set-chassis handset
end script
