# etc/config/includes.chroot/lib/systemd/system/ofono.service.d/pine.conf
[Service]
ExecStartPost=/usr/bin/enable-modem