# etc/config/includes.chroot/lib/udev/rules.d/99-flash.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="white:flash", RUN+="/bin/chmod a+w /sys/class/leds/%k/brightness"
