# etc/config/includes.chroot/lib/systemd/system/modem-start.service
[Unit]
Description=modem power cycle to boot it

[Service]
Type=oneshot
ExecStart=/usr/sbin/modem-start.sh

[Install]
WantedBy=multi-user.target
