# etc/config/includes.chroot/lib/systemd/system/pinetab2-wifi-enable.service
[Unit]
Description=Enable PineTab2 experimental WiFi modules

[Service]
Type=oneshot
ExecStart=/usr/bin/enable-pinetab2-wifi.sh
