# etc/config/includes.chroot/lib/systemd/system/pinetab2-wifi-disable.service
[Unit]
Description=Disable PineTab2 experimental WiFi modules

[Service]
Type=oneshot
ExecStart=/usr/bin/disable-pinetab2-wifi.sh
